# supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:mariadb]
command=/usr/sbin/mysqld
priority=10

[program:proxy]
command=/bin/bash -c "cd /proxy/proxy/build/ && java -jar proxy.jar"
priority=100 